# server/dtek-schedule.service
[Unit]
Description=DTEK Schedule API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/dtek-display/server
Environment="PATH=/home/ubuntu/dtek-display/server/venv/bin"
Environment="API_PASSWORD=API_PASSWORD"
ExecStart=/home/ubuntu/dtek-display/server/venv/bin/gunicorn -w 1 -b 0.0.0.0:5000 --timeout 120 --access-logfile - --error-logfile - server:app
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

